# etc/xdg/autostart/seahorse-daemon.desktop
[Desktop Entry]
Type=Application
Name=Seahorse passphrase caching agent
Exec=/usr/bin/seahorse-daemon
Icon=seahorse
OnlyShowIn=XFCE;
